# tfgm_tram_analyzer/rootfs/etc/services.d/tfgm-tram-analyzer/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "TfGM Tram Analyzer stopped"
