# tfgm_tram_analyzer/rootfs/etc/services.d/tfgm-tram-analyzer/run
#!/usr/bin/with-contenv bashio

# Read add-on options and export as environment variables
export TRAM_WEBSITE_URL=$(bashio::config 'tram_website_url')
export DESTINATIONS=$(bashio::config 'destinations | join(",")')
export SCAN_INTERVAL=$(bashio::config 'scan_interval')
export PYTHONUNBUFFERED=1

bashio::log.info "Starting TfGM Tram Analyzer v3.0"
bashio::log.info "Stop URL: ${TRAM_WEBSITE_URL}"
bashio::log.info "Destinations: ${DESTINATIONS}"
bashio::log.info "Scan interval: ${SCAN_INTERVAL}s"

exec python3 /app/api.py
